# specs/002-prod-broadcast-improvements/deploy/systemd/yt-dlp-update.timer
[Unit]
Description=Run yt-dlp updater weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
