# specs/002-prod-broadcast-improvements/deploy/systemd/yt-dlp-update.service
[Unit]
Description=Update yt-dlp in per-release venv
After=network.target

[Service]
Type=oneshot
User=tgstream
WorkingDirectory=/opt/tg_video_streamer/current
ExecStart=/opt/tg_video_streamer/current/scripts/yt-dlp-update.sh

[Install]
WantedBy=multi-user.target
